# installers/linux/systemd/pantheon-l2.service
[Unit]
Description=PantheonChain L2 Drachma Node (PoS/BFT)
Documentation=https://docs.pantheonchain.org/operations
After=network-online.target
Wants=network-online.target
Requires=network-online.target

[Service]
Type=simple
User=pantheon
Group=pantheon
WorkingDirectory=/var/lib/pantheon/l2

ExecStart=/usr/local/bin/pantheon-l2-drachma-node \
    --layer=l2 \
    --config=/etc/pantheon/mainnet/l2.conf

ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=60
KillMode=process

Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=120

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pantheon-l2

NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/pantheon/l2 /var/log/pantheon
ProtectHome=yes
CapabilityBoundingSet=
AmbientCapabilities=
SecureBits=noroot

LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=8G
CPUQuota=400%

Environment=PANTHEON_NETWORK=mainnet
Environment=PANTHEON_LAYER=l2

[Install]
WantedBy=multi-user.target
